# projects/acpid/service/acpid.service
[Unit]
Description=ACPI Event Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/acpid

[Install]
WantedBy=multi-user.target
