# projects/bitlbee/service/bitlbee@.service
[Unit]
Description=BitlBee Per-Connection Server

[Service]
User=bitlbee
ExecStart=/usr/sbin//bitlbee -I
StandardInput=socket
