# projects/cumin/service/cumin.service
[Unit]
Description=MRG Management Console
After=postgresql.target

[Service]
ExecStart=/usr/bin/cumin --syslog

[Install]
WantedBy=multi-user.target
