# projects/389-ds-base/target/dirsrv.target
[Unit]
Description=389 Directory Server
After=network.target

[Install]
WantedBy=multi-user.target
