# projects/bluez/service/bluetooth.service
[Unit]
Description=Bluetooth Service

[Service]
Type=dbus
BusName=org.bluez
ExecStart=/usr/sbin/bluetoothd -n

[Install]
WantedBy=bluetooth.target
