# projects/cups/service/cups.service
[Unit]
Description=CUPS Printing Service

[Service]
ExecStart=/usr/sbin/cupsd -f

[Install]
Also=cups.socket cups.path
WantedBy=printer.target
