# projects/cloud-init/target/cloud-config.target
[Unit]
Description=Cloud-config availability
Requires=cloud-init-local.service cloud-init.service
