# projects/and/service/and.service
[Unit]
Description=Auto Nice Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/and

[Install]
WantedBy=multi-user.target
