# projects/at/service/atd.service
[Unit]
Description=AT Job Spooling Tools Daemon

[Service]
ExecStart=/usr/sbin/atd -f

[Install]
WantedBy=multi-user.target
