# projects/atop/service/atop.service
[Unit]
Description=advanced interactive monitor

[Service]
ExecStart=/usr/bin/atop

[Install]
WantedBy=multi-user.target
